# deploy/iul-discord-bot.service
[Unit]
Description=IUL Discord Bot
After=network.target iul-backend.service

[Service]
Type=simple
User=iul
WorkingDirectory=/opt/iul-discord-bot
EnvironmentFile=/etc/iul-discord-bot.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iul-discord-bot
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
